# macos_service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.company.desktopagent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/DesktopMonitoringAgent.app/Contents/MacOS/DesktopMonitoringAgent</string>
        <string>--service</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/DesktopAgent.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/DesktopAgent.error.log</string>
    <key>LSUIElement</key>
    <string>1</string>
    <key>ProcessType</key>
    <string>Background</string>
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
